# air-net-led.service
[Unit]
Description=Monitor Internet Access With Airplane LED

[Service]
ExecStart=/usr/bin/air-net-led.sh

[Install]
WantedBy=multi-user.target
